# com.market-data.fetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.market-data.fetch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>--project</string>
        <string>/Users/isdance/Desktop/projects/yahoo-finance-data</string>
        <string>market-data</string>
        <string>fetch</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/isdance/Desktop/projects/yahoo-finance-data</string>

    <!-- Run Mon-Fri at 16:35 ET (UTC-4 = 20:35 UTC, UTC-5 = 21:35 UTC) -->
    <!-- Using 20:35 UTC covers EDT; adjust to 21:35 for EST if needed -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>35</integer>
        </dict>
        <!-- Tuesday -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>35</integer>
        </dict>
        <!-- Wednesday -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>35</integer>
        </dict>
        <!-- Thursday -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>35</integer>
        </dict>
        <!-- Friday -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>35</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/isdance/.market_data/logs/fetch.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/isdance/.market_data/logs/fetch.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
